# wombat.sv.conf
[program:wombat]
directory=/src/Wombat
command=/src/node-v6.2.2-linux-arm64/bin/npm start
